# distrobox/init.service
# This is free and unencumbered software released into the public domain.

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/opt/init.sh
