# distrobox/init.service
[Unit]
Description=Spark service

[Service]
ExecStart=/opt/init.sh

[Install]
WantedBy=multi-user.target
